# pg.louso.ru/etc/systemd/system/proxy-route.service
[Unit]
Description=policy-route for ipset proxied traffic
After=gre-p2p.service netfilter-persistent.service
Requires=gre-p2p.service

[Service]
Type=oneshot
RemainAfterExit=yes

# таблица и правило идемпотентно
ExecStart=/usr/bin/bash -c '\
grep -q "^200[[:space:]]\+proxy" /etc/iproute2/rt_tables || echo "200 proxy" >> /etc/iproute2/rt_tables; \
ip rule add fwmark 0x42 table proxy 2>/dev/null || true; \
ip route add default via 10.24.0.2 dev gre-p2p table proxy 2>/dev/null || true \
'

[Install]
WantedBy=multi-user.target
