# pg.louso.ru/etc/systemd/system/proxy-pg2-route.service
[Unit]
Description=policy-route for ipset pg2_proxy traffic
After=netfilter-persistent.service
Requires=netfilter-persistent.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/usr/sbin/ipset create pg2_proxy hash:ip -exist

# idemponent route rule
ExecStart=/usr/bin/bash -c '\
  grep -q "^201[[:space:]]\+pg2" /etc/iproute2/rt_tables || echo "201 pg2" >> /etc/iproute2/rt_tables; \
  ip rule add fwmark 0x43 table pg2 2>/dev/null || true; \
  ip route replace default via 10.26.0.2 dev gre-pg2 table pg2; \
  ip route replace 10.200.0.0/24 dev wg200 scope link table pg2 \
'

[Install]
WantedBy=multi-user.target 