# pg.gena.host/etc/systemd/system/gre-keepalive.timer
[Unit]
Description=Ping GRE peer every 30 s

[Timer]
OnBootSec=0
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
