# pg.gena.host/etc/systemd/system/gre-pg2-keepalive.timer
[Unit]
Description=Ping GRE pg↔pg2 every 30s

[Timer]
OnBootSec=0
OnUnitActiveSec=30

[Install]
WantedBy=timers.target 